# grub-btrfs-rc
#!/sbin/openrc-run

pidfile="/var/run/grub-btrfs-rc.pid"

depend() {
   use localmount
}

start() {
    ebegin "Starting grub-btrfs daemon"
    start-stop-daemon --background --start --exec /usr/bin/grub-btrfsd
    eend $?
}

